# argus.service
[Unit]
Description=Argus Search API
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/github/argus
ExecStart=/bin/bash /home/ubuntu/github/argus/scripts/start-server.sh
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
SyslogIdentifier=argus
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
